# deploy/vps/euro-core-api.service
[Unit]
Description=Euro Smart - Core API
Documentation=https://github.com/lephuvu1994/euro-smart-server
After=network.target

[Service]
Type=simple
User=node
Group=node
WorkingDirectory=/opt/euro-smart-server
EnvironmentFile=/opt/euro-smart-server/.env
ExecStart=/usr/bin/node dist/apps/core-api/main.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=euro-core-api

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/euro-smart-server/logs

# Resource limits
MemoryMax=512M
CPUQuota=100%

# Graceful shutdown
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
