# deploy/vps/euro-worker-service.service
[Unit]
Description=Euro Smart - Worker Service
Documentation=https://github.com/lephuvu1994/euro-smart-server
After=network.target

[Service]
Type=simple
User=node
Group=node
WorkingDirectory=/opt/euro-smart-server
EnvironmentFile=/opt/euro-smart-server/.env
ExecStart=/usr/bin/node dist/apps/worker-service/main.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=euro-worker-service

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/euro-smart-server/logs

MemoryMax=384M
CPUQuota=50%

KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=60

[Install]
WantedBy=multi-user.target
