# deploy/vps/euro-iot-gateway.service
[Unit]
Description=Euro Smart - IoT Gateway
Documentation=https://github.com/lephuvu1994/euro-smart-server
After=network.target

[Service]
Type=simple
User=node
Group=node
WorkingDirectory=/opt/euro-smart-server
EnvironmentFile=/opt/euro-smart-server/.env
ExecStart=/usr/bin/node dist/apps/iot-gateway/main.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=euro-iot-gateway

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/euro-smart-server/logs

MemoryMax=384M
CPUQuota=50%

KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
